# scripts/macos/com.ziri.listener.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ziri.listener</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/aditya/Desktop/Ziri/.venv/bin/python3</string>
        <string>/Users/aditya/Desktop/Ziri/run_listener.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/aditya/Desktop/Ziri</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/aditya/Library/Logs/Ziri/stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/aditya/Library/Logs/Ziri/stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/aditya/Desktop/Ziri/.venv/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
